# supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=/haproxy-start

[program:sshd]
command=/usr/sbin/sshd -D
